# beanstalk/supervisord.conf
[supervisord]
nodaemon=true

[program:play-beanstalk]
command=/opt/docker/bin/play-beanstalk