# beanstalk/supervisord.conf
[supervisord]
nodaemon=true

[program:play-beanstalk]
command=/opt/docker/bin/play-beanstalk
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0